# factory-reset.target
[Unit]
Description=Factory Reset
Conflicts=shutdown.target systemd-readahead-collect.service systemd-readahead-replay.service
AllowIsolate=yes
