# system-update.target
[Unit]
Description=System Update
Documentation=man:systemd.special(7)
Requires=sysinit.target
After=sysinit.target
