# system-restore.target
[Unit]
Description=System Restore
Conflicts=shutdown.target systemd-readahead-collect.service systemd-readahead-replay.service
AllowIsolate=yes
